# systemd/owntone.service.d/override.conf
[Unit]
Wants=owntone-record_player-input.service
BindsTo=owntone-record_player-input.service
